# bobctf/crontab/ssh.sh
#!/bin/bash
ssh guest@125.131.189.33 -p 10101
